# lib/parsedmarc.timer
[Unit]
Description=Daily run of Parsedmarc

[Timer]
OnBootSec=10min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
